# systemd/minikube.service
# Starts the minikube cluster — and therefore the public site — at boot.
#
# Why this exists: once the site is served from Kubernetes, the compose stack's
# `restart: unless-stopped` policies are no longer what brings it back. Kubernetes restarts the
# workloads on its own, but nothing restarts KUBERNETES: the minikube node is a container inside
# Colima's VM, and neither the VM nor the node comes up by itself. Without this unit, every reboot
# takes project-platform.me offline until someone logs in and types `minikube start`.
#
# A *user* unit, like colima.service, because the Colima VM and the docker socket belong to this user
# ($HOME/.config/colima) and so does the kubeconfig this writes. `loginctl enable-linger` is what lets
# a user unit start at boot with no login session — WITHOUT IT THIS UNIT WILL NOT RUN AT BOOT. See
# the install steps in the platform-orchestration README.
[Unit]
Description=minikube (platform Kubernetes cluster)
# Ordering only — After= does not pull colima in, so Requires= makes the dependency real: with the
# VM down there is no docker daemon for the node container to exist in.
After=colima.service
Requires=colima.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/git-workspace/claude-workspace/platform-orchestration/systemd/minikube-boot.sh
ExecStop=/home/linuxbrew/.linuxbrew/bin/minikube stop
# Booting a QEMU VM and then a Kubernetes control plane inside it is slow; the default 90s timeout
# would kill it mid-start.
TimeoutStartSec=900
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
